# python_job/crontab
* * * * * python /app/scripts/app.py > /dev/stdout
